# docker/mysql-8.0/conf/mysqld.conf
[program:mysqld]
directory=/
command=mysqld
process_name=mysqld
numprocs=1
autostart=true
autorestart=true
user=mysql
group=mysql
redirect_stderr=true
stdout_logfile=/data/logs/mysqld-supervisor.log
